# LaunchAgents/com.fideguch.bochi-healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fideguch.bochi-healthcheck</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/fumito_ideguchi/.claude/scripts/hooks/bochi-s3-healthcheck.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/fumito_ideguchi</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/fumito_ideguchi/.claude/logs/bochi-healthcheck.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/fumito_ideguchi/.claude/logs/bochi-healthcheck.err.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
